# systemd/user/ibus.service
[Unit]
Description=Ibuse Daemon

[Service]
Type=forking
ExecStart=ibus-daemon & disown
Restart=on-failure

[Install]
WantedBy=default.target
